# 1.7.3.2/overlay/etc/update-motd.d/70-prestashop
#!/bin/bash

cat <<EOF

Prestashop is installed in /var/www/html/.
EOF

cat /usr/share/doc/scaleway/prestashop/README 2>/dev/null
